# Common/Supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=java -jar /root/Jenkins/jenkins.war

[program:subversion]
command=svnserve -d --foreground -r /root/DemoSVNServer

